# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:homeautomation]
command=/bin/bash build_and_run.sh %(ENV_ip)s